# pr-downloader.timer
[Unit]
Description=Run pr-downloader every 5min

[Timer]
OnActiveSec=5sec
OnUnitActiveSec=5min

[Install]
WantedBy=timers.target
