# spads_config_updater.timer
[Unit]
Description=Run Spads Config Updater Every Day

[Timer]
OnCalendar=* * 12:00:00

[Install]
WantedBy=timers.target
